# com.vram.eddings-index.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vram.eddings-index</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/VRAM/00-09_System/01_Tools/TheEddingsIndex/.build/release/ei-cli</string>
        <string>watch</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Volumes/VRAM/00-09_System/01_Tools/TheEddingsIndex</string>
    <key>StandardOutPath</key>
    <string>/Users/ronaldeddings/Library/Logs/vram/reality/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ronaldeddings/Library/Logs/vram/reality/error.log</string>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
